# linux/samples/rc.local
#!/bin/bash
echo "Starting ThingPlug Device SDK..."
/home/nvidia/thingplug-device-sdk-C/linux/samples/AutoStartup.sh
